# meta-topic-platform/recipes-graphics/tigervnc/tigervnc-service/x0vncserver.socket
[Socket]
ListenStream=5900

[Install]
WantedBy=sockets.target
